# pkgroot/Library/LaunchAgents/com.rarias84.MarvelComicsStatusBar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.rarias84.MarvelComicsStatusBar</string>
<key>ProgramArguments</key>
<array>
  <string>/Applications/MarvelComicsStatusBar.app/Contents/MacOS/MarvelComicsStatusBar</string>
</array>
<key>RunAtLoad</key><true/>
<key>KeepAlive</key><dict><key>SuccessfulExit</key><false/></dict>
<key>LimitLoadToSessionType</key><string>Aqua</string>
<key>StandardOutPath</key><string>/tmp/MarvelComicsStatusBar.out.log</string>
<key>StandardErrorPath</key><string>/tmp/MarvelComicsStatusBar.err.log</string>
</dict>
</plist>